# lib/templates/sphinx/sphinx_app.initd.centos.erb
#! /bin/sh
# sphinx_<%= application %>: Sphinx search daemon for <%= application %>
#
# chkconfig: - 86 14
# description: sphinx search daemon
# processname: searchd
#
# Author: Gabriel Handford http://ducktyper.com

# Source function library
. /etc/rc.d/init.d/functions

RETVAL=0

DESC="sphinx daemon (<%= application %>)"
NAME=searchd
DAEMON=<%= sphinx_prefix %>/bin/searchd
CONFIGFILE=<%= sphinx_conf_path %>
PIDFILE=<%= sphinx_pid_path %>
INDEX_DIR=<%= sphinx_index_root %>
USER=<%= user %>


start() {
  echo -n $"Starting $DESC: "
  daemon --user $USER --pidfile $PIDFILE $DAEMON --config $CONFIGFILE
  RETVAL=$?
  echo
  return $RETVAL;  
}

stop() {
  echo -n $"Stopping $DESC: "
  killproc -p $PIDFILE
  RETVAL=$?
  echo 
  return $RETVAL;
}

reload() {
  echo -n $"Reloading $DESC: "
  killproc -p $PIDFILE searchd -HUP
  RETVAL=$?
  echo 
  return $RETVAL;
}

clean() {
  echo -n $"Cleaning $DESC: "
  pid=`cat $PIDFILE 2>/dev/null`
  
  if checkpid $pid; then
    failure $"$NAME clean"
  else
    rm -f $INDEX_DIR/*.spl
    RETVAL=$?
    [ "$RETVAL" -eq 0 ] && success $"$NAME clean" || failure $"$NAME clean"    
  fi
  echo
  return $RETVAL;
}

case "$1" in
  start)
    clean
  	start
	;;
  stop)
  	stop
	;;
  reload)
  	reload
  ;;
  restart)
  	stop
  	sleep 1
  	start
	;;
	clean)
    echo -n "Cleaning $DESC"
    clean
    RETVAL=$?;
    ;;
  *)
	  echo "Usage: $0 {start|stop|restart|reload|clean}" >&2
	  RETVAL=3;
	;;
esac

exit $RETVAL;